# hardware/heartbeat_service/heartbeat.service
[Unit]
Description=Raspberry Pi Heartbeat Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=team10
WorkingDirectory=/home/team10/RealMail/hardware/heartbeat_service
ExecStart=/usr/bin/python3 /home/team10/RealMail/hardware/heartbeat_service/heartbeat_service.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=heartbeat-service

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
